# mini-books/Linux Tasks/Week 2/Day 8-9/Linux Chapter 3 (Processes & Services)/my_service.service
[Unit]
Description=My Python Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 /path/to/my_service.py
Restart=always
User=your_user

[Install]
WantedBy=multi-user.target
